# Deploy/Spawn/roles/jobs_deployment_role/files/eventprocessor.conf
[program:eventprocessor]
command=/srv/Tempr/Jobs/eventprocessor.py
user=root
autostart=true
autorestart=true
