# Deploy/Spawn/roles/jobs_deployment_role/files/sendpush.conf
[program:sendpush]
command=/srv/Tempr/Jobs/sendpush.py
user=root
autostart=true
autorestart=true
