# Deploy/Spawn/roles/jobs_deployment_role/files/sendmail.conf
[program:sendmail]
command=/srv/Tempr/Jobs/sendmail.py
user=root
autostart=true
autorestart=true
