# Deploy/Spawn/roles/jobs_deployment_role/files/temps.conf
[program:temps]
command=/srv/Tempr/Jobs/temps.py
user=root
autostart=true
autorestart=true
